# devops/packer/provisioners/file/etc/systemd/system/consul.service
[Unit]
Description="HashiCorp Consul - Service Mesh Made Easy"
Documentation=https://www.consul.io/
Requires=network-online.target
After=network-online.target
ConditionFileNotEmpty=/opt/consul/config/consul.hcl

[Service]
User=consul
Group=consul
ExecStart=/bin/bash /opt/consul/bin/run_consul.sh
ExecReload=/usr/local/bin/consul reload
KillMode=process
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target